# ansible_collections/homelab/proxmox_lxc/roles/loki/templates/loki.service.j2
[Unit]
Description=Loki Log Aggregation System
Documentation=https://grafana.com/docs/loki/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ loki_user }}
Group={{ loki_group }}
ExecStart={{ loki_bin_dir }}/loki -config.file={{ loki_config_dir }}/loki.yml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=loki
KillMode=mixed
KillSignal=SIGINT

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ loki_data_dir }} {{ loki_log_dir }}
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
