# ansible_collections/homelab/proxmox_lxc/roles/openwrt/templates/openwrt.service.j2
[Unit]
Description=OpenWrt Router in LXC
After=network.target

[Service]
Type=forking
User=root
WorkingDirectory=/opt/openwrt
ExecStart=/usr/local/bin/manage-openwrt start
ExecStop=/usr/local/bin/manage-openwrt stop
ExecReload=/usr/local/bin/manage-openwrt restart
RemainAfterExit=yes
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target