# ansible_collections/homelab/proxmox_lxc/roles/pve_exporter/templates/pve-exporter.service.j2
[Unit]
Description=Prometheus PVE Exporter
After=network.target

[Service]
Type=simple
User={{ pve_exporter_user }}
Group={{ pve_exporter_group }}
ExecStart=/opt/pve-exporter/venv/bin/pve_exporter \
  --config.file {{ pve_exporter_config_dir }}/pve.yml \
  --web.listen-address 0.0.0.0:{{ pve_exporter_port }}
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths={{ pve_exporter_data_dir }}

[Install]
WantedBy=multi-user.target
