# ansible_collections/homelab/proxmox_lxc/roles/promtail/templates/promtail.service.j2
[Unit]
Description=Promtail Log Shipper
Documentation=https://grafana.com/docs/loki/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ promtail_user }}
Group={{ promtail_group }}
ExecStart={{ promtail_bin_dir }}/promtail -config.file={{ promtail_config_dir }}/promtail.yml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=promtail
KillMode=mixed
KillSignal=SIGINT

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ promtail_data_dir }} {{ promtail_log_dir }}
ReadOnlyPaths=/var/log
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Allow access to system logs and journal
SupplementaryGroups=adm systemd-journal

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
