# ansible_collections/homelab/proxmox_lxc/roles/media_stack/templates/media-stack.service.j2
# Managed by Ansible - {{ ansible_managed }}
[Unit]
Description=Media Stack (Sonarr, Radarr, Bazarr, Prowlarr, qBittorrent, Jellyfin)
Requires=docker.service
After=docker.service network-online.target mnt-nas.mount
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{ media_stack_data_dir }}
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=on-failure
RestartSec=10
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
