# ansible_collections/homelab/proxmox_lxc/roles/step_ca/templates/step-ca.service.j2
[Unit]
Description=step-ca Local Certificate Authority
Documentation=https://smallstep.com/docs/step-ca
After=network.target

[Service]
User={{ step_ca_user }}
Group={{ step_ca_group }}
Environment=STEPPATH={{ step_ca_config_dir }}
ExecStart=/usr/local/bin/step-ca {{ step_ca_config_dir }}/config/ca.json \
  --password-file {{ step_ca_config_dir }}/password
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ step_ca_config_dir }} {{ step_ca_data_dir }}

[Install]
WantedBy=multi-user.target
