# ansible_collections/homelab/common/roles/monitoring_agent/templates/promtail.service.j2
[Unit]
Description=Promtail - Log collection agent for Loki
Documentation=https://grafana.com/docs/loki/latest/clients/promtail/
After=network.target

[Service]
User={{ monitoring_user }}
Group={{ monitoring_user }}
Type=simple
ExecStart=/usr/local/bin/promtail \
  -config.file={{ monitoring_config_dir }}/promtail.yml \
  -log.level=warn
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{ monitoring_data_dir }} {{ monitoring_log_dir }}
SupplementaryGroups=adm systemd-journal

[Install]
WantedBy=multi-user.target
